# scripts/jarvis-note/com.jarvis.note-daemon.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!--
  launchd agent TEMPLATE for the JARVIS-note daemon's SERVE loop (cloud → you).

  DO NOT install this until the first `serve --allow-writes` run has been done
  supervised and verified clean (the guardrail) — this loop performs live writes to
  the encrypted Raycast DB. It shares the Scheduler's .claude-note-board.lock, so it
  is safe to run alongside the Scheduler's own agents.

  Install (after a clean supervised run):
    1. Replace every __PLACEHOLDER__ below (absolute paths + the operator secret).
    2. cp to ~/Library/LaunchAgents/com.jarvis.note-daemon.plist
    3. launchctl load -w ~/Library/LaunchAgents/com.jarvis.note-daemon.plist
  Logs land in /tmp; check them after load. Unload: launchctl unload -w <path>.
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.jarvis.note-daemon</string>

    <key>ProgramArguments</key>
    <array>
        <string>__PYTHON3__</string>           <!-- e.g. /usr/bin/python3 -->
        <string>__JARVIS_REPO__/scripts/jarvis-note/daemon.py</string>
        <string>serve</string>
        <string>--allow-writes</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>JARVIS_APP_URL</key>
        <string>__JARVIS_APP_URL__</string>     <!-- e.g. https://mydearestjarvis.vercel.app -->
        <key>RAYCAST_INGEST_SECRET</key>
        <string>__RAYCAST_INGEST_SECRET__</string>
        <key>SCHEDULER_REPO</key>
        <string>__SCHEDULER_REPO__</string>      <!-- e.g. /Users/david/Developer/Claude - Scheduler -->
    </dict>

    <key>RunAtLoad</key>
    <true/>
    <!-- KeepAlive restarts the long-poll loop if it ever exits. -->
    <key>KeepAlive</key>
    <true/>
    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>StandardOutPath</key>
    <string>/tmp/com.jarvis.note-daemon.out.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/com.jarvis.note-daemon.err.log</string>
</dict>
</plist>
